# app/share/e2guardian/scripts/e2guardian.service
[Unit]
Description=E2guardian Web filtering
After=network.target

[Service]
TasksMax=infinity
LimitNOFILE=65535
Type=forking
IgnoreSIGPIPE=no
GuessMainPID=no
ExecStart=/app/sbin/e2guardian
ExecReload=/app/sbin/e2guardian -r
UMask=027

[Install]
WantedBy=multi-user.target
